# 99-xpad_joystick_disabler.rules
KERNEL=="js*", RUN+="/usr/local/bin/xpad_joystick_disabler"
